# clock-mcu/src/firmware/device/etc/init.d/coprocessor
#!/bin/sh /etc/rc.common

START=80

USE_PROCD=1
BIN="/usr/bin/mcu/update.sh"


boot() {
    $BIN
}
